# ansible/templates/rivetr.service.j2
# /etc/systemd/system/rivetr.service
# Generated by Ansible — do not edit by hand.

[Unit]
Description=Rivetr PaaS Engine
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory={{ rivetr_install_dir }}
ExecStart={{ rivetr_install_dir }}/rivetr --config {{ rivetr_install_dir }}/rivetr.toml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rivetr

[Install]
WantedBy=multi-user.target
